# RealSync-AI-Prototype/deploy/realsync-ai.service
[Unit]
Description=RealSync AI Inference Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/RealSync-AI-Prototype
Environment=PORT=5100
Environment=HOST=0.0.0.0
ExecStart=/usr/bin/python3 -m serve.app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
